# common/profile.d/000-sysutils
#!/bin/bash

case "$(uname)" in
	Darwin)
		ensureCmd lf "brew install lf"

		setOrInstall JENV_HOME \
			"$HOME/.jenv" \
			"brew install jenv" 
		;;
	Linux)
		setOrInstall JENV_HOME \
			"$HOME/.jenv" \
			"git clone https://github.com/gcuisinier/jenv.git $JENV_HOME" 

		eval "$( jenv init - )"

		;;
esac

JAVA_HOME="$( jenv javahome )"
export JAVA_HOME
export JRE_HOME="${JAVA_HOME}/jre"
export MAVEN_OPTS="-XX:+UseG1GC -XX:MaxHeapFreeRatio=40"

# LF -- The dopest filemanager
LFCD_HOME="$HOME/.config/lf/"
LFCD="${LFCD_HOME}/lfcd.sh"
if [ -f "$LFCD" ]; then
	source "$LFCD"
fi
