# common/profile.d/000-appspath
#!/bin/sh

for app in $(ls $APPS); do
	export PATH=$APPS/$app/bin:$PATH
done
