# common/profile.d/100-optpath
#!/bin/sh

export JAVA_HOME=$OPT/java
export JETTY_HOME=$OPT/jetty

export PATH=$JAVA_HOME/bin:$JETTY_HOME/bin:$PATH
