# common/profile.d/100-homepath
#!/bin/sh

case `uname` in
  Darwin)
	[ -e ~/.fonts/] && source ~/.fonts/*.sh
  ;;
  Linux)
	echo ""
  ;;
esac
