# common/profile.d/102-shellevals
#!/bin/sh

eval `dircolors ~/.dircolors`
