# common/profile.d/102-shellevals
#!/bin/sh
[[ -e $(which gdircolors) ]] && eval $(gdircolors ~/.dircolors)
[[ -e $(which dircolors) ]] && eval $(dircolors ~/.dircolors)
