# common/profile.d/203-shellvars-bb
export SSO_JWT_EXTERNAL_EXPIRATION="36000"
export SSO_JWT_EXTERNAL_RENEW="144000"
export SSO_JWT_EXTERNAL_NOT_VALID_AFTER="57600"
