# common/profile.d/502-shellalias-bp
#!/bin/sh

alias dbp="dev cd business-platform"
alias dshop="dev cd shopify"
